# tests/cron_to_json_test_happy.cron
# Intro comment

# First action comment
0 * * * * echo 'hello world'
# Second action
# comment
1 * * * * echo 'farewell'
# Outro